# gatherer_role/templates/fahrpreis_gatherer.service.j2
[Unit]
Description=collect fahrpreise

[Service]
ExecStart=/opt/fahrpreis_gatherer/run.sh
User=fahrpreis_gatherer
