# gatherer_role/templates/fahrpreis_gatherer.timer.j2
[Unit]
Description=schedule fahrpreis_gatherer

[Timer]
OnCalendar=*-*-* *:00/3:00
Persistent=true

[Install]
WantedBy=timers.target
